# package/stockdb-backup.service
[Unit]
Description=Backup StockDB

[Service]
Type=oneshot
ExecStart=/bin/bash /var/backups/backup-stockdb.sh
